# blog_portafolio/service-conf/gunicorn-blog.socket
Description=Blog Gunicorn Socket

[Socket]
ListenStream=/run/gunicorn-blog.sock
User=www-data
Mode=600

[Install]
WantedBy=sockets.target